# service/rtorrent/run
#!/bin/sh

mkdir -p /data
TERM=xterm rtorrent -n -o import=/conf/rtorrent.rc
